# infra/ec2/systemd/brain-core.service
[Unit]
Description=Brain Core FastAPI orchestrator
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/var/brain/apps/brain-core
Environment=PATH=/home/ubuntu/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=BRAIN_VAULT_PATH=/var/brain
Environment=BRAIN_DB_PATH=/var/brain/apps/brain-core/db/brain.sqlite
ExecStart=/home/ubuntu/.local/bin/uv run python -m brain_core
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
