# infra/ec2/systemd/claude-creds-sync.timer
[Unit]
Description=Periodic Claude credentials sync to AWS Secrets Manager

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s
Unit=claude-creds-sync.service

[Install]
WantedBy=timers.target
